# rebooter.service
[Unit]
Description=Rebooter
After=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/reboot
