# rebooter.timer
[Unit]
Description=Reboot timer

[Timer]
OnCalendar=06:00
Unit=rebooter.service

[Install]
WantedBy=timers.target
